# systemd/clutch-server.service
[Unit]
Description=OpenClutch Next.js Server
Documentation=https://github.com/dbachelder/clutch
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/dan/src/clutch
EnvironmentFile=/home/dan/src/clutch/.env.local
Environment=NODE_ENV=production
Environment=PORT=3002
ExecStart=/home/dan/.volta/tools/image/node/22.22.0/bin/node ./node_modules/next/dist/bin/next start -p 3002
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clutch-server

[Install]
WantedBy=default.target
